# infra/remote-work/images/remote-workspace/scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_RESOLUTION)s -ac +extension GLX +render -noreset
user=worker
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xfce]
command=/bin/bash -lc "export HOME=/home/worker; export DISPLAY=%(ENV_DISPLAY)s; dbus-launch --exit-with-session startxfce4"
user=worker
autorestart=true
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
command=/bin/bash -lc "x11vnc -display %(ENV_DISPLAY)s -rfbport %(ENV_VNC_PORT)s -rfbauth /home/worker/.vnc/passwd -forever -shared -noxdamage"
user=worker
autorestart=true
startsecs=5
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:websockify]
command=/bin/bash -lc "websockify --web /opt/novnc %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s"
user=worker
autorestart=true
startsecs=3
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
